# deployment/supervisor.conf
; ============================================================
; omarhesham.school — Supervisor Configuration
; File: /etc/supervisor/conf.d/omarhesham-worker.conf
;
; Apply with:
;   sudo supervisorctl reread
;   sudo supervisorctl update
;   sudo supervisorctl start all
; ============================================================

; ── Default queue (general jobs, badge checks, etc.) ────────
[program:omarhesham-default]
process_name=%(program_name)s_%(process_num)02d
command=php /home/forge/omarhesham.school/current/artisan queue:work redis --queue=default --sleep=3 --tries=3 --max-time=3600
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=forge
numprocs=2
redirect_stderr=true
stdout_logfile=/home/forge/.forge/omarhesham-worker.log
stopwaitsecs=3600

; ── Email queue (consent emails, badge notifications, etc.) ─
[program:omarhesham-emails]
process_name=%(program_name)s_%(process_num)02d
command=php /home/forge/omarhesham.school/current/artisan queue:work redis --queue=emails --sleep=3 --tries=5 --max-time=3600
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=forge
numprocs=1
redirect_stderr=true
stdout_logfile=/home/forge/.forge/omarhesham-emails.log
stopwaitsecs=3600
